# install/rc.local
#!/bin/bash
cd /home/pi/betabot

sleep 10
git pull

./start 2> /home/pi/betabot/error.log > /home/pi/betabot/log.log &

