# install/rc.local
#!/bin/bash

cd /home/pi/betabot
./start 2> /home/pi/betabot/error.log > /home/pi/betabot/log.log &
